# scripts/lxc/villaparis.service
[Unit]
Description=Villa Paris Gestionale
Documentation=https://github.com/Dis-Astro/villa-paris-gestionale
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/villaparis
EnvironmentFile=/opt/villaparis/.env
ExecStart=/usr/bin/yarn start
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/villaparis

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=villaparis

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
